# bots/bug-watcher/systemd/localbitcoiners-bug-watcher.service
[Unit]
Description=localbitcoiners — poll relay.mynostr.app for new bug reports and create GitHub issues
Documentation=file://%h/localbitcoiners/bots/bug-watcher/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# ⚠️ Adjust these two paths to wherever the localbitcoiners repo lives on
# the server (mirrors mynostr's unit, which uses its own repo path).
WorkingDirectory=%h/localbitcoiners
ExecStart=%h/localbitcoiners/bots/bug-watcher/run.sh
# Bound runtime — a stuck WS connection should never stall the timer.
TimeoutStartSec=60s
Nice=10
StandardOutput=journal
StandardError=journal
